# home/greg/.config/systemd/user/code.service
[Unit]
Description=code

[Service]
WorkingDirectory=%h/code/arch
ExecStart=%h/code/arch/code.sh

