# home/greg/.config/systemd/user/dropbox-fix.service
[Unit]
Description=dropbox-fix

[Service]
Type=simple
ExecStart=/usr/bin/sh ~/code/arch/dropbox-fix.sh

[Install]
WantedBy=multi-user.target
