# home/greg/.config/systemd/user/fetch.service
[Unit]
Description=fetch
After=graphical-session.target

[Service]
Type=simple
EnvironmentFile=%h/code/arch/fetch.env
PassEnvironment=GITHUB_PAT_APSIS GITHUB_PAT_EFFICY GITHUB_PAT_GREG
ExecStart=%h/code/arch/fetch.py

[Install]
WantedBy=graphical-session.target

