# home/greg/.config/systemd/user/history.service
[Unit]
Description=history

[Service]
WorkingDirectory=%h/code/history
ExecStart=%h/code/history/sync.sh

