# home/greg/.config/systemd/user/notes.service
[Unit]
Description=notes

[Service]
WorkingDirectory=%h/code/notes
ExecStart=%h/code/notes/sync.sh

