# home/greg/.config/systemd/user/nvidia.service
[Unit]
Description=nvidia

[Service]
Type=oneshot
ExecStart=%h/code/arch/nvidia.sh
