# home/greg/.config/systemd/user/nvidia.timer
[Unit]
Description=nvidia

[Timer]
OnBootSec=15

[Install]
WantedBy=timers.target
