# home/greg/.config/systemd/user/passwords.service
[Unit]
Description=passwords

[Service]
WorkingDirectory=%h/code/passwords
ExecStart=%h/code/passwords/sync.sh

