# home/greg/.config/systemd/user/passwords.timer
[Unit]
Description=passwords

[Timer]
OnCalendar=15:00:00
Persistent=true

[Install]
WantedBy=timers.target
