# home/greg/.config/systemd/user/passwords.timer
[Unit]
Description=passwords
After=network-online.target

[Timer]
OnCalendar=00/1:00
Persistent=true

[Install]
WantedBy=timers.target
