# home/greg/.config/systemd/user/sync-periodic.service
[Unit]
Description=sync-periodic

[Service]
Type=oneshot
RemainAfterExit=false
WorkingDirectory=%h/code/arch
ExecStart=%h/code/arch/sync.sh
