# home/greg/.config/systemd/user/wall.service
[Unit]
Description=wall

[Service]
Type=oneshot
RemainAfterExit=false
WorkingDirectory=%h/code/arch
ExecStart=%h/code/arch/wall.sh
